# user/display-setup.service
[Unit]
Before=awesome.service

[Service]
ExecStart=/home/fixje/bin/display-setup.sh
Type=oneshot

[Install]
WantedBy=default.target
