# user/ipython-notebook.service
[Service]
ExecStart=/home/fixje/hacks/notebooks/start.sh
Type=forking

[Install]
WantedBy=default.target
